# deploy/pyhn_uwsgi.conf
description "pyhn uwsgi"
start on runlevel [2345]
stop on runlevel [06]

exec /home/pyhn/PyHackerNews/venv/bin/uwsgi --ini /home/pyhn/PyHackerNews/deploy/uwsgi.ini
